# cfg/msb-mqtt.service
[Unit]
Description=Motion Sensor Box MQTT Service
After=msb-broker.service
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=msb-mqtt
Type=simple
Restart=always
RestartSec=1
User=pi
ExecStart=/usr/bin/python3 /home/pi/motion-sensor-box/src/mqtt-client-msb/src/msb_mqtt.py

[Install]
WantedBy=multi-user.target